# power-monitor.service
[Unit]
Description=HeadUnit Power Monitor
After=multi-user.target
Before=headunit.service

[Service]
Type=simple
User=mike
ExecStart=/usr/bin/python3 /home/mike/HeadUnit/power_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
